# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; 1. Virtual display
[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x720x24 -ac
autorestart=true
priority=1
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; 2. VNC server - streams the virtual display
[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -listen 127.0.0.1 -forever -shared -quiet
autorestart=true
priority=2
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; 3. websockify - bridges VNC to WebSocket so browsers can connect
;    serves noVNC static files from /usr/share/novnc
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 127.0.0.1:5900
autorestart=true
priority=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; 4. Chromium browser
[program:chrome]
command=/usr/bin/chromium
    --no-sandbox
    --disable-dev-shm-usage
    --display=:1
    --window-size=1280,720
    --window-position=0,0
    --disable-gpu
    --disable-infobars
    --disable-extensions
    --noerrdialogs
    --disable-session-crashed-bubble
    --no-first-run
    --user-data-dir=/tmp/chrome-profile
    --start-maximized
    %(ENV_START_URL)s
environment=DISPLAY=":1"
autorestart=false
priority=4
stdout_logfile=/dev/null
stderr_logfile=/dev/null
